# castoro-pgctl/etc/init.d/cstartd.svccfg.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='application/castoro/cstartd' type='service' version='0'>
    <create_default_instance enabled='true'/>
    <single_instance/>
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    <dependency name='filesystem-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>
    <dependency name='multi-user' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/multi-user'/>
    </dependency>
    <dependency name='configuration_files' grouping='require_all' restart_on='restart' type='path'>
      <service_fmri value='file://localhost/etc/castoro/pgctl.conf'/>
      <service_fmri value='file://localhost/etc/castoro/peer.conf'/>
    </dependency>
    <dependency name='executables' grouping='require_all' restart_on='restart' type='path'>
      <service_fmri value='file://localhost/usr/local/bin/cstartd'/>
    </dependency>
    <exec_method name='start' type='method' exec='/etc/init.d/cstartd start' timeout_seconds='60'/>
    <exec_method name='stop' type='method' exec='/etc/init.d/cstartd stop' timeout_seconds='60'/>
    <stability value='Unstable'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>Castoro Start Daemon</loctext>
      </common_name>
      <documentation>
        <doc_link name='Castoro' uri='https://github.com/ToshiyukiTerashita/castoro'/>
      </documentation>
    </template>
  </service>
</service_bundle>
